# systemd/vcc-ready.target
[Unit]
Description=VCC cluster ready target
Requires=vcc-net.service