# systemd/vcc-services.target
[Unit]
Description=VCC launch service units
Requires=vcc-wait4deps.service
After=vcc-wait4deps.service
Before=vcc-register.service