# systemd/vcc-wait4deps.service
[Unit]
Description=VCC wait for dependent services on other nodes
PartOf=vcc.service
After=vcc-keys.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/node /vccjs/wait4deps.js

[Install]
RequiredBy=vcc.service

